# shelly-firmware-updater/rootfs/etc/services.d/shelly/run
#!/bin/sh
set -e
# Gunicorn in der venv starten (Pfad im Dockerfile gesetzt)
exec /opt/venv/bin/gunicorn -b 0.0.0.0:8099 app:app --workers 4 --threads 2 --timeout 90
