# shelly-firmware-updater/rootfs/etc/services.d/shelly/run
#!/usr/bin/with-contenv bashio
set -e
# Fester Ingress-Port (wie in config.json)
exec gunicorn -b 0.0.0.0:8099 app:app --workers 4 --threads 2 --timeout 90
