# shelly-firmware-updater/rootfs/etc/services.d/shelly/run
#!/usr/bin/with-contenv bashio
set -e

# Standard-Ingress-Port (aus config.json). Falls Supervisor ENV setzt, das nutzen:
PORT="${INGRESS_PORT:-8099}"

# Gunicorn aus der venv starten
exec gunicorn -b 0.0.0.0:${PORT} app:app --workers 4 --threads 2 --timeout 90
